# openwrt/etc/init.d/aoa-proxy-forward
#!/bin/sh /etc/rc.common

# TODO: this is really just a placeholder for now, so that I remember the general setup in OpenWRT.