# systemd/aoa-proxy-announce@.service
[Unit]
Description="announce our presence to a USB connected android device"

[Service]
ExecStart=/usr/lib/aoa-proxy-announce %I
