# openwrt/etc/hotplug.d/usb/99-aoa-proxy
#!/bin/sh /etc/rc.common


case "${ACTION}" in
    bind)
        case "${PRODUCT}" in
            18d1/2d00/* | 18d1/2d01/* | 18d1/2d04/* | 18d1/2d05/*)
                service aoa-proxy-forward start
                ;;
            *)
                . /etc/os-release
                MANUFACTURER="aoa-proxy"
                MODEL="$(echo -e "${OPENWRT_DEVICE_MANUFACTURER} ${OPENWRT_DEVICE_PRODUCT}")"
                VERSION="${PRETTY_NAME}"
                SERIAL="$(uci get system.@system[0].hostname)"
                DESCRIPTION="""$MODEL
IP: $(ip -o route | egrep -v "(^default|dev (virbr|docker))" | sed -r 's/.* dev (\S*) .+ src (\S+) .*/\2\t\1/')
"""
                URL="https://github.com/jo-bitsch/aoa-proxy/"
                MODEL_AND_SERVICES="""$MODEL
ssh"""
                /usr/sbin/aoa-proxy \
                    --port="${DEVICENAME}" \
                    --announce \
                    --manufacturer="$MANUFACTURER" \
                    --model="$MODEL_AND_SERVICES" \
                    --model-version="$VERSION" \
                    --serial="$SERIAL" \
                    --url="https://github.com/jo-bitsch/aoa-proxy/" \
                    --description="$DESCRIPTION"
                ;;
        esac
        ;;
    unbind)
        service aoa-proxy-forward stop "${DEVICENAME}"
        ;;
esac